# update.crontab
0 */2 * * * ~/vcs/git/dotfiles/update.sh >> ~/tmp/logs/dotfiles.log
